# packages/sway/update-wallpaper.timer
[Unit]
Description=Update wallpaper every 5 minutes

[Timer]
OnBootSec=0s
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
