# packages/sway/update-wallpaper.service
[Unit]
Description=Update desktop wallpaper

[Service]
ExecStart=%h/.local/bin/update-wallpaper

[Install]
WantedBy=default.target
