# packages/i3/i3-update-wallpaper.timer
[Unit]
Description=Update i3 window manager wallpaper every 5 minutes

[Timer]
OnBootSec=0s
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
