# packages/i3/i3-update-wallpaper.service
[Unit]
Description=Update i3 window manager wallpaper

[Service]
ExecStart=%h/.local/bin/i3-update-wallpaper

[Install]
WantedBy=default.target
